# infra/oracle/pocketdock-gateway.service
[Unit]
Description=PocketDock Android control gateway
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/pocketdock/gateway
EnvironmentFile=/etc/pocketdock/gateway.env
Environment=HOME=/tmp
ExecStart=/usr/local/bin/node /opt/pocketdock/gateway/src/server.mjs
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
