# infra/oracle/pocketdock-quick-tunnel.service
[Unit]
Description=PocketDock temporary Cloudflare Quick Tunnel
After=network-online.target pocketdock-gateway.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate --url http://127.0.0.1:8787
Restart=no
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
